# bin/apps.service
[Unit]
Description=Apps

[Service]
ExecStart=/home/michael/xbns/bin/apps
KillMode=mixed
Type=simple
Restart=always

[Install]
WantedBy=multi-user.target