# bin/apps.service
[Unit]
Description=Apps

[Service]
ExecStart=/home/michael/xbns/bin/apps
KillMode=mixed
Type=idle
Restart=no

[Install]
WantedBy=multi-user.target