# bin/ota.service
[Unit]
Description=Ota

[Service]
ExecStart=/home/michael/xbns/bin/ota
KillMode=mixed
Type=simple
Restart=always

[Install]
WantedBy=multi-user.target