# bin/pong.service
[Unit]
Description=Pong

[Service]
ExecStart=/home/michael/xbns/bin/pong
KillMode=mixed
Type=simple
Restart=always

[Install]
WantedBy=multi-user.target