# bin/pong.service
[Unit]
Description=Pong

[Service]
ExecStart=/home/michael/xbns/bin/pong
KillMode=mixed
Type=idle
Restart=no

[Install]
WantedBy=multi-user.target