# bin/xbns.service
[Unit]
Description=XBNS

[Service]
ExecStart=/home/michael/xbns/bin/xbns
KillMode=mixed
Type=idle
Restart=no

[Install]
WantedBy=multi-user.target